# system/hotspot/hotspot.service
[Unit]
Description=WiFi Hotspot for Jetson Field Kit
After=network-online.target NetworkManager.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/home/box/jetson-orin-nano-field-kit/system/hotspot/hotspot.conf
ExecStart=/home/box/jetson-orin-nano-field-kit/system/hotspot/setup-hotspot.sh
ExecStop=/home/box/jetson-orin-nano-field-kit/system/hotspot/setup-hotspot.sh --stop
Restart=on-failure
RestartSec=10

# Wait for WiFi to connect AND USB adapter to be available
ExecStartPre=/bin/bash -c 'for i in {1..60}; do if nmcli -t -f STATE g 2>/dev/null | grep -q connected && ip link show wlxccbabd28d984 &>/dev/null; then exit 0; fi; sleep 2; done; echo "Timeout waiting for network/adapter"; exit 0'

[Install]
WantedBy=multi-user.target
